# hiveden-dashboard.service
[Unit]
Description=Hiveden Dashboard
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/hiveden-dashboard
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000
# EnvironmentFile=-/etc/default/hiveden-dashboard

[Install]
WantedBy=multi-user.target
